# deploy/agenthook.service
# systemd user unit for the agenthook control-plane.
#
# Install (as the dedicated `agenthook` user):
#   agenthook install-service --port 8080 --write   # generates this under ~/.config/systemd/user/
#   # then EDIT ExecStart to bind loopback (see below), and:
#   systemctl --user daemon-reload
#   systemctl --user enable --now agenthook
#
# Requires: `loginctl enable-linger agenthook` so the unit runs without a login
# session, and the `agenthook` user in the `docker` group.
#
# This file is a reference/template — `install-service --write` emits an
# equivalent; the one change you MUST make is binding 127.0.0.1 (not 0.0.0.0) so
# only the co-located reverse proxy can reach it and X-Forwarded-For is honored.

[Unit]
Description=agenthook control-plane (webhook runner)
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
# Bind loopback: the TLS reverse proxy (nginx/caddy) is the only public ingress.
ExecStart=%h/.local/bin/agenthook serve --host 127.0.0.1 --port 8080
Environment=AGENTHOOK_HOME=%h/.agenthook
Restart=on-failure
RestartSec=3
# Hardening (optional, tighten as your host allows):
NoNewPrivileges=true
ProtectSystem=full
PrivateTmp=true

[Install]
WantedBy=default.target
